# supplemental/init_d/chatd.in
#!/bin/sh

### BEGIN INIT INFO
# Provides:		chatd
# Required-Start:	$network $remote_fs $syslog
# Required-Stop:	$network $remote_fs $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	Connme chat server 
# Description:		Connme chat server listening on 7681
### END INIT INFO

#log_daemon_msg "1)" "${NAME}"

set -e

DAEMON="${CHAT_INSTALL_DIR}/chat-server"
CHATD_OPTS="--daemonize -s"
NAME="chatd"
PATH="/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin"

test -x "${DAEMON}" || exit 0

# log_daemon_msg "1)" "${NAME}"

#. /lib/lsb/init-functions

# log_daemon_msg "2)" "${NAME}"


case "${1}" in
	start)
#		log_daemon_msg "Starting chat server" "${NAME}"

		echo "Starting chat server"
		start-stop-daemon --start --oknodo --exec ${DAEMON} -- ${CHATD_OPTS}
		
		if ! ps -C chat-server | grep -qs "chat-server"
		then
#			log_warning_msg "chat-server failed - probably already run."
			echo "chat-server failed - probably already run."
			exit 1
		fi

#		log_end_msg 0
		;;

	stop)
#		log_daemon_msg "Stopping chat server" "${NAME}"

		echo "Stopping chat server"
#		start-stop-daemon --stop --oknodo --exec ${DAEMON}
		kill `more /tmp/.lwsts-lock`

#		log_end_msg 0
		;;

	reload|force-reload|restart)
		${0} stop
		${0} start
		;;

	status)
		status_of_proc "${DAEMON}" "chat server"
		;;

	*)
		echo "Usage: ${0} {start|stop|restart|reload|status}"
		exit 1
		;;
esac

exit 0
